# dns-analyzer.service
[Unit]
Description=DNS Server Analyzer - Distributed DNS Measurement
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=oneshot
User=root
Group=root
WorkingDirectory=/root/dns-analyzer
ExecStart=/bin/bash /root/dns-analyzer/run.sh

# Resource limits
TimeoutStartSec=0
TimeoutStopSec=60s

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=dns-analyzer

# Security
PrivateTmp=true

[Install]
WantedBy=multi-user.target
