# dns-analyzer.service
[Unit]
Description=DNS Server Analyzer - Distributed DNS Measurement
After=network-online.target postgresql.service
Wants=network-online.target
Documentation=

[Service]
Type=oneshot
User=root
Group=root
WorkingDirectory=/root/dns-analyzer

# Use run.sh which handles everything (venv, dependencies, execution)
ExecStart=/bin/bash /root/dns-analyzer/run.sh

# Restart on failure (but not on success since Type=oneshot)
Restart=on-failure
RestartSec=5min

# Resource limits
TimeoutStartSec=0
TimeoutStopSec=60s

# Logging configuration
StandardOutput=journal
StandardError=journal
SyslogIdentifier=dns-analyzer

# No security hardening when running as root
PrivateTmp=true

[Install]
WantedBy=multi-user.target
