# dns-analyzer.timer
[Unit]
Description=DNS Analyzer Timer - Runs every hour after completion
Requires=dns-analyzer.service
Documentation=https://github.com/your-repo/dns-analyzer

[Timer]
# First run: 5 minutes after system boot
OnBootSec=5min

# Subsequent runs: 1 hour AFTER service completes
OnUnitActiveSec=1h

# Store last activation time on disk (survives reboots)
Persistent=true

# Allow ±1 minute randomization to prevent thundering herd
# (useful if you run this on multiple servers)
AccuracySec=1min

# Specify which service to activate
Unit=dns-analyzer.service

[Install]
WantedBy=timers.target
